# systemd/muslimtify.service
[Unit]
Description=Prayer Time Notification Check
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/muslimtify check
StandardOutput=journal
StandardError=journal
